# dot_config/systemd/user/chezmoi-auto-sync.service
[Unit]
Description=Sync chezmoi source and push dotfiles
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/chezmoi-auto-sync

[Install]
WantedBy=default.target
